# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:uwsgi]
command = /usr/bin/uwsgi --ini /home/facelook/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
